# deploy/emg-gateway.service
# A+7: runs the EMG gateway as an always-on service on the Raspberry Pi. Install to /etc/systemd/system/.
[Unit]
Description=EMG Telemetry gateway (notch filter + AWS cloud uplink)
After=network-online.target mosquitto.service
Wants=network-online.target

[Service]
User=pi
WorkingDirectory=/home/pi/emg-gateway
# load the EMG char-device driver first (+ runs it as root, - ignores "already loaded")
ExecStartPre=-+/usr/sbin/insmod /home/pi/emg-driver/emgactivity.ko
ExecStart=/home/pi/emg-gateway/venv/bin/python /home/pi/emg-gateway/emg_gateway.py
Restart=on-failure
RestartSec=3
Environment=PYTHONUNBUFFERED=1

[Install]
WantedBy=multi-user.target
